# collector/systemd/observal-collector.service
[Unit]
Description=Observal AV Collector Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=observal
Group=observal
EnvironmentFile=/etc/observal/observal.env
ExecStart=/opt/observal/venv/bin/observal-collector
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=observal-collector

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/observal
PrivateTmp=true

[Install]
WantedBy=multi-user.target
